# packages/backend/config/redis/redis.service
[Unit]
Description=Redis for Ionfi
After=syslog.target
[Service]
Type=simple
ExecStart=/usr/local/bin/redis-server /etc/redis/redis.conf
RestartSec=5s
Restart=always
[Install]
WantedBy=multi-user.target
